# stacks/tailscale/tailscale-health.service
[Unit]
Description=Tailscale healthcheck (oneshot)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tailscale-healthcheck.sh

